# modules/video/video.service
[Unit]
Description=ffs-video
After=multi-user.target

[Service]
Type=simple
Restart=on-failure
User=pi
WorkingDirectory=/home/pi/ffs/modules/video
ExecStart=/home/pi/ffs/modules/video/venv/bin/python -m src.main

[Install]
WantedBy=multi-user.target